# video-converter.service
[Unit]
Description=P4Frame Video Converter Service
After=multi-user.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/P4Frame

# Main service command with config file
ExecStart=/usr/bin/python3 /home/pi/P4Frame/video_converter.py --config /home/pi/P4Frame/p4frame_raspi4.conf

# Environment variables
Environment="HOME=/home/pi"

# Restart policy
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=600

# Resource limits
MemoryMax=500M
MemorySwapMax=500M
CPUQuota=180%
Nice=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=p4frame-converter

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/share/study/lgy_photoframe/data /home/pi/P4Frame/logs /var/log/video_converter /home/pi/P4Frame
NoNewPrivileges=true

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target